# etc/services.d/postfix/run
#!/bin/bash
set -e
postfix check
exec postfix start-fg